# printer_app.service
[Unit]
Description=Printer Application (Flask + ESC/POS Agent)
After=network.target

[Service]
Type=simple
User=printer-app
ExecStart=/opt/printer_application/.venv/bin/printer-app
Restart=always
RestartSec=5
Environment=PRINTER_APP_DATA=/var/lib/printer_app
# Uncomment to override defaults:
# Environment=PRINTER_APP_PORT=5000
# Environment=PRINTER_APP_HOST=0.0.0.0
StandardOutput=journal
StandardError=journal
SyslogIdentifier=printer-app

[Install]
WantedBy=multi-user.target
